# synth.service
[Unit]
Description=Synth - Self-hosted AI research engine
Documentation=https://github.com/Bentlybro/synth
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/bently/clawd/projects/synth
Environment="PATH=/home/bently/.cargo/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/bently/clawd/projects/synth/.env
ExecStart=/home/bently/clawd/projects/synth/target/release/synth
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
